# curator4-1754350331.conf
[program:curator4-1754350331]
command=python3 /home/curator4/curator4-1754350331/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1754350331/test.err.log
stdout_logfile=/home/curator4/curator4-1754350331/test.out.log